# deploy/telegram-bot.service.example
# Example systemd unit for the Telegram bot (VPS deployment).
# Copy to /etc/systemd/system/telegram-bot.service and adjust paths/user.
# Then: sudo systemctl daemon-reload && sudo systemctl enable telegram-bot && sudo systemctl start telegram-bot

[Unit]
Description=Telegram Lead Tracking & Signal Mirroring Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/telegram-bot
EnvironmentFile=/opt/telegram-bot/.env
ExecStart=/opt/telegram-bot/venv/bin/gunicorn app.main:app -w 1 -k uvicorn.workers.UvicornWorker -b 0.0.0.0:8000
# If PORT is set in .env, use: -b 0.0.0.0:${PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
